# docker/xori-api.conf
[program:xori-api]
directory=/root/xori/gui
command=node src/server.js
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
